# usr/lib/systemd/system/ima-offload.target
[Unit]
Description=Immutarch Offload

[Install]
WantedBy=local-fs.target
